# recipes-openharmony/openharmony/files/appspawn.service
[Unit]
Description=OpenHarmony Appspawn service

[Service]
Type=oneshot
ExecStart=/bin/echo "OpenHarmony Appspawn Service"
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
