# recipes-openharmony/openharmony/files/distributedsched.service
[Unit]
Description=OpenHarmony Distributed Schedule service
Requires=appspawn.service
After=appspawn.service

[Service]
Type=exec
# should run as system:system
ExecStart=/usr/bin/sa_main /system/profile/distributedsched.xml

[Install]
WantedBy=multi-user.target
