# recipes-openharmony/openharmony/files/openharmony-preinit.service
[Unit]
Description=OpenHarmony pre-init setup

[Service]
Type=oneshot
ExecStart=/bin/sh /usr/sbin/openharmony-preinit

[Install]
WantedBy=multi-user.target
