# recipes-openharmony/openharmony/files/param.service
[Unit]
Description=OpenHarmony Param Service

[Service]
Type=notify
ExecStartPre=install -m 755 -d /run/openharmony/param
ExecStartPre=install -m 755 -d /data/param
ExecStart=/usr/bin/param_service

[Install]
WantedBy=multi-user.target
