# recipes-openharmony/openharmony/files/param.service
[Unit]
Description=OpenHarmony Param Service

[Service]
Type=oneshot
ExecStart=/bin/echo "OpenHarmony Param Service"
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
