# recipes-openharmony/openharmony/files/tokensync.service
[Unit]
Description=OpenHarmony TokenSyncManager
Requires=samgr.service

[Service]
Type=exec
# should run as system:system
ExecStart=/usr/bin/sa_main /system/profile/token_sync_service.xml
Restart=on-failure
StartLimitInterval=300
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
